# optional/automation/com.example.daily-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.daily-summary</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/CHANGE_ME/.claude/hooks/daily-summary.sh</string>
    </array>

    <!-- End of day: 18:00 local. Companion to com.example.daily-plan (08:00). -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/CHANGE_ME/.claude/cache/daily-summary.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/CHANGE_ME/.claude/cache/daily-summary.log</string>
</dict>
</plist>
